# rootfs/systemd/camera-hal.service
[Unit]
Description=Ubuntu GSI Camera HAL — V4L2/libcamera HIDL Bridge
Documentation=https://github.com/zerofrip/Ubuntu_Touch_HIDL_GSI
After=hwbinder-bridge.service systemd-udevd.service
Wants=hwbinder-bridge.service
ConditionPathExists=/usr/lib/ubuntu-gsi/hidl/camera/camera_hal.sh

[Service]
Type=simple
ExecStart=/usr/lib/ubuntu-gsi/hidl/camera/camera_hal.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
Environment=LIBCAMERA_LOG_LEVELS=*:WARN

# Security hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/ubuntu-gsi /data/uhl_overlay /dev
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
