# rootfs/systemd/gnss-hal.service
[Unit]
Description=Ubuntu GSI GNSS HAL — gpsd/Geoclue HIDL Bridge
Documentation=https://github.com/zerofrip/Ubuntu_Touch_HIDL_GSI
After=hwbinder-bridge.service systemd-udevd.service
Wants=hwbinder-bridge.service
ConditionPathExists=/usr/lib/ubuntu-gsi/hidl/gnss/gnss_hal.sh

[Service]
Type=simple
ExecStart=/usr/lib/ubuntu-gsi/hidl/gnss/gnss_hal.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

# Security hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/ubuntu-gsi /data/uhl_overlay /dev
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
