# rootfs/systemd/input-hal.service
[Unit]
Description=Ubuntu GSI Input HAL — Touchscreen/Input Classifier (HIDL synthetic)
Documentation=https://github.com/zerofrip/Ubuntu_Touch_HIDL_GSI
After=hwbinder-bridge.service systemd-udevd.service
Wants=hwbinder-bridge.service
Before=lomiri.service

[Service]
Type=simple
ExecStart=/usr/lib/ubuntu-gsi/hidl/input/input_hal.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

# Security hardening
ProtectSystem=strict
ReadWritePaths=/data /run/ubuntu-gsi /tmp /dev/input /etc/libinput /usr/local/bin
NoNewPrivileges=false
PrivateTmp=no

[Install]
WantedBy=graphical.target
