# rootfs/systemd/ubuntu-gsi-firstboot.service
[Unit]
Description=Ubuntu GSI First Boot Setup (HIDL)
Documentation=https://github.com/zerofrip/Ubuntu_Touch_HIDL_GSI
ConditionPathExists=!/data/uhl_overlay/.firstboot_complete
After=local-fs.target
Before=hwbinder-bridge.service lomiri.service

[Service]
Type=oneshot
ExecStart=/usr/lib/ubuntu-gsi/firstboot.sh
RemainAfterExit=yes
TimeoutStartSec=0

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ubuntu-gsi-firstboot

[Install]
WantedBy=multi-user.target
