# rootfs/systemd/ubuntu-gsi-setup-wizard.service
[Unit]
Description=Ubuntu GSI GUI Setup Wizard
Documentation=https://github.com/zerofrip/Ubuntu_Touch_GSI
ConditionPathExists=/data/uhl_overlay/.setup_wizard_pending
ConditionPathExists=!/data/uhl_overlay/.setup_wizard_complete
After=lomiri.service graphical.target dbus.service

[Service]
Type=oneshot
User=ubuntu
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0
Environment=DISPLAY=:0
ExecStartPre=/bin/sleep 5
ExecStart=/usr/lib/ubuntu-gsi/setup-wizard.sh
RemainAfterExit=yes
TimeoutStartSec=0

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ubuntu-gsi-setup-wizard

[Install]
WantedBy=graphical.target
