# rootfs/systemd/vibrator-hal.service
[Unit]
Description=Ubuntu GSI Vibrator HAL — sysfs/FF Vibration HIDL Bridge
Documentation=https://github.com/zerofrip/Ubuntu_Touch_HIDL_GSI
After=hwbinder-bridge.service systemd-udevd.service
Wants=hwbinder-bridge.service
ConditionPathExists=/usr/lib/ubuntu-gsi/hidl/vibrator/vibrator_hal.sh

[Service]
Type=simple
ExecStart=/usr/lib/ubuntu-gsi/hidl/vibrator/vibrator_hal.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

# Security hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/ubuntu-gsi /data/uhl_overlay /dev /sys/class/timed_output /sys/class/leds /usr/local/bin
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
