# deprecated/systemd-units-old/face-hal.service
[Unit]
Description=Ubuntu GSI Face Auth HAL — Howdy/IR Camera HIDL Bridge
Documentation=https://github.com/zerofrip/Ubuntu_Touch_HIDL_GSI
After=hwbinder-bridge.service systemd-udevd.service camera-hal.service
Wants=hwbinder-bridge.service
ConditionPathExists=/usr/lib/ubuntu-gsi/hidl/face/face_hal.sh

[Service]
Type=simple
ExecStart=/usr/lib/ubuntu-gsi/hidl/face/face_hal.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

# Security hardening
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/ubuntu-gsi /data/uhl_overlay /dev /sys /etc/pam.d /etc/howdy
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
