# deprecated/systemd-units-old/wifi-hal.service
[Unit]
Description=Ubuntu GSI WiFi HAL — Vendor WiFi HIDL Bridge
Documentation=https://github.com/zerofrip/Ubuntu_Touch_HIDL_GSI
After=hwbinder-bridge.service network-pre.target
Wants=hwbinder-bridge.service
Before=NetworkManager.service wpa_supplicant.service

[Service]
Type=simple
ExecStart=/usr/lib/ubuntu-gsi/hidl/wifi/wifi_hal.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=60

# Security hardening
ProtectSystem=strict
ReadWritePaths=/data /run/ubuntu-gsi /tmp /etc/wpa_supplicant /etc/NetworkManager /lib/firmware /run/wpa_supplicant
NoNewPrivileges=false
PrivateTmp=no

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wifi-hal

[Install]
WantedBy=multi-user.target
